# root/etc/services.d/curseforgebot/run
#!/usr/bin/with-contenv bash

cd /app/curseforgebot

exec s6-setuidgid abc /usr/lib/jvm/jdk-15/java -jar CurseForge2Discord.jar
