# root/etc/services.d/curseforgebot/run
#!/usr/bin/with-contenv bash

cd /config

exec s6-setuidgid abc /usr/lib/jvm/jdk-15/bin/java -jar /app/curseforgebot/CurseForge2Discord.jar
