# root/etc/services.d/curseforgebot/run
#!/usr/bin/with-contenv bash

cd /app/curseforgebot

exec \
	s6-setuidgid abc java -jar curseforgebot.jar
